# supervisor/supervisord-san-884.conf
[supervisord]
logfile = /tmp/supervisord-san-pca.log

[program:run-experiments]
directory=.
command=python run-study-noniid-mnist.py --opt --log --nep-log --not-pure --attackers-num 8 --output-prefix "pca_niid_np_moe_atk_SAN_884"
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

