# supervisor/supervisord-san-492.conf
[supervisord]
logfile = /tmp/supervisord-san-492.log

[program:run-experiments]
directory=.
command=python run-study-iid-mnist.py --avg --log --nep-log --attackers-num 75 --output-prefix "avg_atk1_iid_SAN_492"
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

