# supervisor/supervisord-san-882.conf
[supervisord]
logfile = /tmp/supervisord-san-pca.log

[program:run-experiments]
directory=.
command=python run-study-noniid-mnist.py --opt --log --nep-log --attackers-num 15 --output-prefix "pca_niid_moe_atk_SAN_882"
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

