# supervisor/supervisord-san-545.conf
[supervisord]
logfile = /tmp/supervisord-san-545.log

[program:run-experiments]
directory=.
command=python run-study-noniid-mnist.py --opt --log --nep-log --not-pure --attackers-num 25 --output-prefix "moe_atk1_niid_npure_SAN_545"
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

