# supervisor/supervisord-san-551.conf
[supervisord]
logfile = /tmp/supervisord-san-551.log

[program:run-experiments]
directory=.
command=python run-study-iid-mnist.py --opt --log --nep-log --not-pure --attackers-num 75 --output-prefix "moe_atk1_iid_npure_SAN_551"
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

