# supervisor/supervisord-san-461.conf
[supervisord]
logfile = /tmp/supervisord-run-461.log

[program:run-experiments]
directory=.
command=python run-study-noniid-mnist.py --avg --log --nep-log --attackers-num 50 --output-prefix "avg_atk1_niid_SAN_461"
autostart=true
autorestart=unexpected
exitcodes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

